# telegram-tts-bot.service
[Unit]
Description=Telegram Text-to-Speech Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mg/telegram-tts-bot
EnvironmentFile=/home/mg/telegram-tts-bot/.env
ExecStart=/home/mg/telegram-tts-bot/.venv/bin/python /home/mg/telegram-tts-bot/bot.py
Restart=always
RestartSec=5
User=mg

[Install]
WantedBy=multi-user.target
